# etc/rc.local
注册开机启动到系统的本文件，即可开机启动。

开机启动顺序上，本文件是最后执行的
